# templates/puppet.cron.erb
30 * * * *    `which puppet` agent --onetime --no-daemonize --verbose
