# templates/puppet.cron.erb
30 * * * *    puppet    puppet agent --onetime --no-daemonize --verbose > /dev/null 2>&1
